# samples/hello.ini
[program:hello]
command=/Library/Java/JavaVirtualMachines/jdk1.8.0_25.jdk/Contents/Home/bin/java -jar /data/jfpal_workspace/architect/samples/hello.jar
process_name=%(program_name)s
directory=/data/jfpal_workspace/architect/samples/
autostart=true
autorestart=true
stderr_logfile=/data/jfpal_workspace/architect/samples/logs/%(program_name)s.err.log
stdout_logfile=/data/jfpal_workspace/architect/samples/logs/%(program_name)s.out.log